# services/conn-monitor.service
[Unit]
Description=TCP Connection Monitor Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/conn-monitor
Restart=always
RestartSec=10
StartLimitInterval=0
Environment="MAIL_SERVER="
Environment="MAIL_PORT=587"
Environment="MAIL_USERNAME="
Environment="MAIL_PASSWORD="
Environment="MAIL_TO="

# 日志配置
SyslogIdentifier=conn-monitor
ReadWritePaths=/var/log

# 安全配置
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes

[Install]
WantedBy=multi-user.target